# rootfs/etc/s6-overlay/s6-rc.d/caddy/run
#!/usr/bin/with-contenv sh

CONFIG_FILE=/config/Caddyfile

# change chdir
cd /config

# run with limited permissions
exec s6-setuidgid $PUID:$PGID /app/caddy run --config=$CONFIG_FILE
